# scripts/crontab
@reboot tmux new -s fidelio
00 12 * * * /home/pi/fidelio-capture-bit/startup.sh
